# smart_lock.conf
[program:smart_lock]
command=bash -c "IDM=YOUR_IDM /home/pi/Repo/smart_lock/app.py" ;
numprocs=1           ;
redirect_stderr=true ;
stdout_logfile=/home/pi/Repo/smart_lock/app.log ;
user=root ;
